# quang.pro/etc/cron.d/restart_html
@reboot   root    /usr/bin/docker restart html